# roles/vector-role/templates/vector.service.j2
# =============================================================================
# Vector systemd unit (generated by Ansible)
# Файл: roles/vector/templates/vector.service.j2
# =============================================================================
[Unit]
Description=Vector Log Agent
Documentation=https://vector.dev/docs/
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{ vector_user }}
Group={{ vector_group }}
ExecStart={{ vector_install_dir }}/bin/vector --config {{ vector_config_file }}
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536

# 🔐 Security Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
PrivateTmp=true
ReadWritePaths={{ vector_log_dir }} {{ vector_config_dir }} {{ vector_data_dir }}

# 🔑 Переменные окружения для паролей
Environment="CLICKHOUSE_PASSWORD=${CLICKHOUSE_PASSWORD:-{{ vector_clickhouse_password }}}"

[Install]
WantedBy=multi-user.target
